# service/lockerd.service
[Unit]
Description=Locker Daemon (lockerd)
After=network.target

[Service]
Type=simple
ExecStart=/opt/locker/bin/lockerd
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
